# data/system76-firmware-daemon.service
[Unit]
Description=System76 Firmware Daemon

[Service]
ExecStart=/usr/bin/system76-firmware-daemon
Restart=no

[Install]
WantedBy=multi-user.target
